# debian/continuum-update.service
[Unit]
Description=Continuum Update Service
Documentation=https://github.com/Jordonh18/continuum
After=network-online.target
Wants=network-online.target

# Don't start automatically - only when triggered
RefuseManualStop=no

[Service]
Type=oneshot

# Run as root since we need to install packages
User=root
Group=root

# Read update parameters from the request file
ExecStart=/opt/continuum/scripts/update_runner.sh

# Clean up on success or failure
ExecStopPost=/bin/rm -f /var/lib/continuum/update_request.json

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=continuum-update

# Security hardening
ProtectSystem=false
ProtectHome=true
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
