# debian/postinst
#!/bin/sh
# postinst script for printer-proxy
#
# See: dh_installdeb(1)

set -e

APP_USER="printer-proxy"
APP_GROUP="printer-proxy"
APP_DIR="/opt/printer-proxy"
CONFIG_DIR="/etc/printer-proxy"
DATA_DIR="/var/lib/printer-proxy"
LOG_DIR="/var/log/printer-proxy"

case "$1" in
    configure)
        # Create system user and group if they don't exist
        if ! getent group "$APP_GROUP" > /dev/null 2>&1; then
            addgroup --system "$APP_GROUP"
        fi
        if ! getent passwd "$APP_USER" > /dev/null 2>&1; then
            adduser --system --ingroup "$APP_GROUP" --home "$APP_DIR" \
                    --no-create-home --gecos "Printer Proxy" "$APP_USER"
        fi

        # Create required directories
        install -d -m 0755 -o root -g "$APP_GROUP" "$CONFIG_DIR"
        install -d -m 0750 -o "$APP_USER" -g "$APP_GROUP" "$DATA_DIR"
        install -d -m 0750 -o "$APP_USER" -g "$APP_GROUP" "$LOG_DIR"

        # Set ownership of application directory
        chown -R root:root "$APP_DIR"
        chown -R "$APP_USER:$APP_GROUP" "$APP_DIR/venv" 2>/dev/null || true

        # Create Python virtual environment
        if [ ! -d "$APP_DIR/venv" ]; then
            echo "Setting up Python virtual environment..."
            python3 -m venv "$APP_DIR/venv"
        fi
        
        # Install Python dependencies
        "$APP_DIR/venv/bin/pip" install --quiet --upgrade pip 2>/dev/null || true
        "$APP_DIR/venv/bin/pip" install --quiet -r "$APP_DIR/requirements.txt" 2>/dev/null || \
            echo "Warning: Some Python packages may not have installed correctly"
        
        chown -R "$APP_USER:$APP_GROUP" "$APP_DIR/venv"

        # Install sudoers configuration for network helper
        install -m 0440 /dev/stdin /etc/sudoers.d/printer-proxy << 'EOF'
# Printer Proxy - allow web app to run privileged helper scripts
Defaults:printer-proxy !requiretty
printer-proxy ALL=(root) NOPASSWD: /opt/printer-proxy/scripts/network_helper.sh
printer-proxy ALL=(root) NOPASSWD: /usr/bin/systemctl start printer-proxy-update.service
printer-proxy ALL=(root) NOPASSWD: /opt/printer-proxy/scripts/update_helper.sh
EOF

        # Install the update service
        if [ -f "$APP_DIR/debian/printer-proxy-update.service" ]; then
            install -m 0644 "$APP_DIR/debian/printer-proxy-update.service" \
                    /etc/systemd/system/printer-proxy-update.service
        fi

        # Generate self-signed SSL certificate if it doesn't exist
        if [ ! -f /etc/ssl/certs/printer-proxy.crt ]; then
            echo "Generating self-signed SSL certificate..."
            CERT_CN=$(hostname -f 2>/dev/null || echo "localhost")
            openssl req -x509 -nodes -days 3650 -newkey rsa:2048 \
                -keyout /etc/ssl/private/printer-proxy.key \
                -out /etc/ssl/certs/printer-proxy.crt \
                -subj "/CN=${CERT_CN}" 2>/dev/null || \
                echo "Warning: Could not generate SSL certificate"
            chmod 600 /etc/ssl/private/printer-proxy.key 2>/dev/null || true
        fi

        # Configure nginx if available
        if [ -d /etc/nginx/sites-available ] && [ -d /etc/nginx/sites-enabled ]; then
            cat > /etc/nginx/sites-available/printer-proxy << 'NGINX_EOF'
server {
    listen 80;
    server_name _;
    return 301 https://$host$request_uri;
}

server {
    listen 443 ssl http2;
    server_name _;

    ssl_certificate /etc/ssl/certs/printer-proxy.crt;
    ssl_certificate_key /etc/ssl/private/printer-proxy.key;
    ssl_protocols TLSv1.2 TLSv1.3;
    ssl_ciphers ECDHE-ECDSA-AES128-GCM-SHA256:ECDHE-RSA-AES128-GCM-SHA256:ECDHE-ECDSA-AES256-GCM-SHA384:ECDHE-RSA-AES256-GCM-SHA384;
    ssl_prefer_server_ciphers off;

    add_header X-Frame-Options "SAMEORIGIN" always;
    add_header X-Content-Type-Options "nosniff" always;

    access_log /var/log/nginx/printer-proxy-access.log;
    error_log /var/log/nginx/printer-proxy-error.log;

    # Serve React frontend
    location / {
        root /opt/printer-proxy/frontend/dist;
        try_files $uri /index.html;
        expires 1h;
        add_header Cache-Control "public, must-revalidate";
    }

    # Proxy API requests to Flask backend
    location /api/ {
        proxy_pass http://127.0.0.1:8080;
        proxy_set_header Host $host;
        proxy_set_header X-Real-IP $remote_addr;
        proxy_set_header X-Forwarded-For $proxy_add_x_forwarded_for;
        proxy_set_header X-Forwarded-Proto $scheme;
        proxy_connect_timeout 30s;
        proxy_send_timeout 30s;
        proxy_read_timeout 30s;
    }

    # Static assets with long cache
    location /static/ {
        alias /opt/printer-proxy/static/;
        expires 1y;
        add_header Cache-Control "public, immutable";
    }

    # Frontend assets with cache
    location /assets/ {
        root /opt/printer-proxy/frontend/dist;
        expires 1y;
        add_header Cache-Control "public, immutable";
    }
}
NGINX_EOF
            ln -sf /etc/nginx/sites-available/printer-proxy /etc/nginx/sites-enabled/
            rm -f /etc/nginx/sites-enabled/default 2>/dev/null || true
            
            if nginx -t 2>/dev/null; then
                if systemctl is-active --quiet nginx 2>/dev/null; then
                    systemctl reload nginx 2>/dev/null || true
                fi
            fi
        fi

        # Reload systemd and enable service
        systemctl daemon-reload
        systemctl enable printer-proxy.service 2>/dev/null || true
        systemctl start printer-proxy.service 2>/dev/null || true

        echo ""
        echo "=============================================="
        echo "  Printer Proxy installed successfully!"
        echo "=============================================="
        echo ""
        HOST_IP=$(hostname -I 2>/dev/null | awk '{print $1}')
        if [ -n "$HOST_IP" ]; then
            echo "Access the web interface at: https://${HOST_IP}"
        else
            echo "Access the web interface via your system's IP address"
        fi
        echo ""
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument '$1'" >&2
        exit 1
        ;;
esac

exit 0
