# debian/printer-proxy.service
[Unit]
Description=Printer Proxy
Documentation=https://github.com/jordonh18/printer-proxy
After=network.target

[Service]
Type=simple
User=printer-proxy
Group=printer-proxy
WorkingDirectory=/opt/printer-proxy
Environment="PATH=/opt/printer-proxy/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/opt/printer-proxy/venv/bin/gunicorn \
    --bind 127.0.0.1:8080 \
    --workers 2 \
    --timeout 30 \
    --access-logfile /var/log/printer-proxy/access.log \
    --error-logfile /var/log/printer-proxy/error.log \
    "app:create_app()"

PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/printer-proxy /var/lib/printer-proxy

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
